# .config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
Restart=always
ExecStart=/bin/sh -c "xremap %h/.config/xremap.yml --watch"
