# .config/autostart/copyq.desktop
[Desktop Entry]
Name=CopyQ
Type=Application
Exec=env QT_QPA_PLATFORM=xcb copyq
Terminal=false
