# .config/autostart/fcitx5.desktop
[Desktop Entry]
Name=Fcitx 5
Type=Application
Exec=/usr/bin/fcitx5
Terminal=false
